# services/Cultibot-Picture.timer
[Unit]
Description=Cultibot Picture

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
